# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:falkordb]
command=redis-server --loadmodule /var/lib/falkordb/bin/falkordb.so --unixsocket /tmp/falkordb.sock --unixsocketperm 777 --save 60 1 --dir /data/graph
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:memkit-api]
command=/usr/local/bin/memkit-api serve --pack /data/memory-pack --host 0.0.0.0 --port 4242
autostart=true
autorestart=true
priority=2
environment=FALKORDB_SOCKET="/tmp/falkordb.sock",LANCEDB_PATH="/data/lance"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
